# templates/cron.hourly
#!/bin/bash
DIR=/var/run/unattended-charm
mkdir -p $DIR
/usr/lib/update-notifier/apt-check >& $DIR/apt-check

